# freebsd/rc.d/mvtserver
#!/bin/sh
# PROVIDE: mvtserver
# REQUIRE: DAEMON NETWORKING
# KEYWORD: shutdown

. /etc/rc.subr

name="mvtserver"
rcvar="mvtserver_enable"

command="/usr/local/bin/mvt-server"
command_chdir="/usr/local/etc/mvt-server"
procname="/usr/local/bin/mvt-server"
pidfile="/var/run/mvtserver.pid"

# Export variables from .env to the environment
if [ -f "/usr/local/etc/mvt-server/.env" ]; then
    export $(grep -v '^#' /usr/local/etc/mvt-server/.env | xargs)
fi

# Load configuration from /etc/rc.conf
load_rc_config $name

# Default values
: ${mvtserver_enable:="NO"}
: ${mvtserver_config:="/usr/local/etc/mvt-server"}
: ${mvtserver_cache:="/var/cache/mvt-server"}
: ${mvtserver_assets:="/usr/local/etc/mvt-server/assets"}
: ${mvtserver_flags:=""}

# Build arguments
command_args="--config ${mvtserver_config} --cache ${mvtserver_cache} --mapassets ${mvtserver_assets} ${mvtserver_flags} >> /var/log/mvt-server.log 2>&1 &"

# Required directories
required_dirs="${mvtserver_config} ${mvtserver_cache} ${mvtserver_assets}"

run_rc_command "$1"

